# systemd/vc-pipeline-v2.timer
[Unit]
Description=Run ESPN Clean v2 pipeline hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
