# systemd/vc-resolver-v2.service.sample
[Unit]
Description=ESPN Clean v2 - Resolver API (FastAPI/uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brad
Group=brad
WorkingDirectory=%h/Projects/ESPN_clean_v2
Environment=PATH=%h/Projects/ESPN_clean_v2/.venv/bin:/usr/bin:/bin
# uses env or config.ini at runtime (no hardcoded IPs)
ExecStart=/usr/bin/env bash -lc 'uvicorn bin.vc_resolver:app --host 0.0.0.0 --port 8094'
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/espnvc-v2/resolver.log
StandardError=append:/var/log/espnvc-v2/resolver.log

ProtectHome=read-only
ProtectSystem=full
ReadWritePaths=%h/Projects/ESPN_clean_v2 /var/log/espnvc-v2
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
