# systemd/legacy/vc-plan.service
[Unit]
Description=ESPN Clean v2.0 - Build Plan + Render Outputs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Projects/ESPN_clean
ExecStart=/usr/bin/python3 bin/build_plan.py --db data/eplus_vc.sqlite3 --valid-hours=72 --tz=America/New_York --note="timer run"
ExecStart=/usr/bin/python3 bin/xmltv_from_plan.py --db data/eplus_vc.sqlite3 --out out/virtual_channels.xml 
ExecStart=/usr/bin/python3 bin/m3u_from_plan.py --db data/eplus_vc.sqlite3 --out out/virtual_channels.m3u --resolver-base=http://192.168.86.72:8093 --cc-host=192.168.86.72 --cc-port=5589 --only-live
