# systemd/legacy/vc-plan.timer
[Unit]
Description=Run ESPN Clean v2.0 plan build hourly

[Timer]
OnCalendar=*:00/60
Persistent=true

[Install]
WantedBy=timers.target
